# rootfs/setup/redis.conf
[program:redis]
command=/bin/redis-server --port 6379
priority=500
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
